# pfsense/crontabs
@reboot root nice -n 19 sh /root/pfsense/ssh-check.sh &
*/5 * * * * root sh /root/pfsense/cpu-check.sh &
0 0 * * * root sh /root/pfsense/disk_usage.sh &